# root/etc/profile.d/editor.sh
#!/bin/sh
# /etc/profile.d/vim
# Sets up vim settings.

# Default text editor setting
EDITOR="vim"
VISUAL="vim"

export VISUAL
export EDITOR

# End of file
