# root/etc/profile.d/miku.sh
#!/bin/sh
# Author: Todd Gaunt
# Print a happy miku!
echo ' 
       ⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⣀⡀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⣀⣀⡀⠤
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣾⣿⣿⣿⣦⡰⣩⣭⡭⠶⢶⠶⠶⠒⠊⠉⠁⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣀⣀⣀⣀⡀⠀⠀⠀⠀⣿⣿⣿⢿⣿⡷⠋⠀⡠⠒⣀⣠⣤⣤⣶⣶⣶⣶⣖⠒⠒
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣠⣴⣶⡲⠯⠛⠋⠀⠀⠀⠀⠀⠀⠉⠁⠂⠒⣿⣿⣿⠀⢹⠇⣰⣿⡾⠿⠛⠛⠛⠛⠋⠙⠙⠛⠛⠻⠲
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣤⡀⠀⠀⠀⠀⢀⡠⠔⠂⢊⣩⢴⡿⣿⣿⣿⣦⣀⠀⠉⠛⢮⣅⡀⠀⠀⠒⠠⢄⣿⣿⣿⢲⣼⣰⠟⠁⠀⣀⣀⣤⣤⣤⣶⣶⣶⣶⣦⣤⣤
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣿⣿⣿⣿⣦⡠⠖⠋⠁⢀⣴⣾⢏⢴⣿⣇⠈⢻⣿⣿⣿⣧⡀⠀⠀⠙⢿⣯⡪⣿⣷⣾⣿⣿⣿⣌⢻⠇⣴⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣀⣀⣀⣠⣿⣿⡿⠟⣁⣠⣔⢒⣿⣿⣿⣿⢫⣿⣿⡞⡄⠀⠙⢿⣿⣿⣿⣦⡀⠀⠀⠙⢷⣄⠙⢟⡛⢿⣿⣿⣿⣿⣿⣿⠿⢿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⠔⠋⠀⢀⡀⢀⣸⣟⣫⢾⣯⣿⣿⣟⣽⣿⣿⡿⢡⡇⣿⣿⢳⡘⢄⠀⠀⠙⢿⣿⣿⣿⣦⠀⠀⠀⠙⠣⡀⠙⢦⣽⣿⣿⣿⣿⣿⣿⣶⡤⣉⡛⠿⣿⣦⣍⠉⠉⠉⠉
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⣔⣥⠖⠋⠉⢉⠝⣩⣽⣿⣿⣀⣽⣿⢫⣾⣿⣿⡿⠁⣼⡇⣿⣿⠈⣧⡈⠳⣄⠀⠀⠙⢿⣿⣿⣷⣄⠀⠀⠀⠙⢦⡀⠙⢿⣿⣿⣿⣿⣿⣿⡇⠈⠉⠲⢤⣉⠻⢿⣶⣄⣀
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⣴⣿⠟⠁⠀⢀⠔⠁⣼⣯⣿⣿⣿⣿⣿⢣⣿⣿⣿⣿⠇⠀⣿⡇⣿⣿⠀⢿⡓⢄⠘⢦⡀⠀⠀⠙⢿⣿⣿⣆⡀⠀⠠⠀⠱⣄⠀⠙⢿⣿⣿⣿⣿⣧⠀⠀⠀⠀⡈⠑⠦⣍⠻⢿
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣠⣿⡿⠋⠀⢀⡴⠋⣠⡞⣽⣿⣿⣿⣿⣿⡷⣿⣿⣿⣿⡟⠀⣼⠃⢧⣿⣿⠀⢸⡇⠀⠳⣄⠹⣦⣄⣤⡬⠿⣿⣿⣿⣦⡀⠀⠀⠈⢳⡀⠀⠙⠿⣿⣿⣿⠀⠀⠀⢠⣟⠀⠀⠈⠙⢶
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣴⣿⠟⠁⠀⡠⢊⣴⣾⣿⢰⣿⡿⠛⠙⢿⣿⣻⣿⣿⣿⣿⠇⢀⠇⠀⠹⣸⣿⡀⢸⡇⠀⠀⡨⠟⢿⣿⣷⣄⠀⠈⠛⢿⣿⣿⣦⡀⠀⠐⠙⢷⣶⣤⣤⣉⣙⣶⣶⠶⠟⠁⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⣾⣿⠏⠀⢀⢔⣵⣿⣿⣿⡧⢘⠁⠀⣀⣴⣿⣿⣿⣟⣹⣿⣿⠀⣞⠀⠀⠀⣿⣿⡇⠀⡿⠀⠀⠀⠀⠀⠙⠛⠿⣷⣄⠀⠀⠉⠻⣿⣿⣦⡀⠀⠀⠙⢿⣿⣿⢹⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢠⣷⡿⠃⠀⡠⣫⣿⣿⣿⠟⠋⠄⢸⠴⠚⣽⣿⣿⣇⣿⣿⣿⣿⣿⣠⡿⠛⠋⠁⢪⢹⣿⠀⢸⠀⠀⠀⠀⠀⠀⣠⣴⣶⣿⣷⣦⡝⠢⣄⠉⠻⢿⣶⣄⠀⠀⠈⠻⣇⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⠀⠀⢠⣻⡿⠁⢀⢞⣵⣿⣿⠟⠁⢀⡴⠚⢄⢇⣼⣿⡿⠟⡇⣿⣿⣿⣿⣿⡸⣧⠀⠀⠀⠀⠁⣧⠁⠀⢃⠀⠀⠀⠀⣼⠛⢡⣿⣿⣿⣷⡝⣷⡰⣍⠓⣤⣝⠻⣷⣄⠀⠈⠲⣕⠦⡀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⠀⢰⣿⡿⠁⣠⢫⣾⣿⡿⢃⡤⠞⠁⠀⠀⢈⣿⣿⡋⠀⠀⢡⣿⣿⣿⣿⡿⡇⢿⠀⠀⣀⣐⠀⠘⣷⣄⠀⠀⠀⠀⠀⠁⠀⠘⠋⢻⣿⣿⣿⠈⣹⡸⣿⣾⠻⣷⣌⢻⣷⣄⠀⠸⡗⢌⢢⠀⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⠀⢠⣳⡿⢁⡔⣵⣿⣿⣫⠖⠋⠀⠀⠀⠀⢀⡾⠋⠀⠀⠀⠀⢸⢹⣿⣿⣿⡇⣿⡘⡀⣾⠟⣻⣿⣄⠘⣎⠳⣌⠀⠀⠀⠀⠀⠀⠲⣝⠋⡹⠋⠀⢹⡇⣿⣿⣷⣌⠻⣆⣿⣿⣷⡄⢱⠀⠡⣇⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⠀⢰⣿⡿⢡⢞⣾⣿⣿⠟⠁⠀⠀⠀⠀⢀⡴⠋⠀⠀⠀⠀⠀⠀⠈⣿⣿⣿⣿⣷⢸⣿⣜⣿⠌⠿⠿⣿⣆⠀⠁⠈⠉⠋⠒⠂⠀⠀⠀⠀⠀⠀⠀⠀⠨⢵⣿⣿⣿⣿⣷⡜⣿⣿⣿⢿⣸⠀⠀⠏⠀⠀⠀⠀⠀
⠀⠀⠀⠀⠀⢠⣿⣿⣱⢣⣿⡿⠛⠁⠀⠀⠀⠀⢀⡴⠋⠀⠀⠀⠀⠀⢀⣤⣶⣾⣿⣻⣿⣿⣿⣿⣿⣿⡌⠃⠀⠰⣋⡻⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢹⣿⡿⠿⠟⠛⣌⡅⠠⠤⠏⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⠀⠀⢠⣿⣿⡗⣵⡿⠋⠀⠀⠀⠀⠀⢀⠴⠋⠀⠀⠀⠀⠀⣠⣶⣿⣿⣿⣿⠟⠹⣹⡏⣧⢿⣿⣿⣿⣆⠀⠀⠂⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⢀⢮⠎⠀⠀⠀⠀⢹⢴⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⠀⢀⣾⣿⡟⣽⠋⠀⠀⠀⠀⠀⢀⠔⠁⠀⠀⠀⠀⢀⣴⣿⣿⣿⣽⣿⣿⡇⡄⠈⣄⢿⡘⢞⣿⣿⣿⣿⡆⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣠⣤⡄⠀⠀⠀⠀⠀⠀⠈⣿⠀⠀⠀⠀⠀⢠⡘⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⢀⣾⣿⣿⠞⠁⠀⠀⠀⠀⢀⡴⠃⠀⠀⠀⠀⢀⣴⣿⣿⣿⣿⣿⣿⣿⣿⣧⣰⠊⠀⠙⢷⣈⠻⣿⣿⣿⣿⡄⠀⠀⠀⠀⢰⣶⣾⣿⡿⠿⠿⡇⠀⠀⠀⠀⡀⠂⠀⡿⠀⠀⠀⠀⠀⠜⠁⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⠀⣼⣿⡿⠃⠀⠀⠀⠀⢀⣴⠟⠀⠀⠀⠀⢀⣴⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣧⢇⠀⠀⠀⠙⠻⠿⢿⣿⣿⣟⢦⡀⠀⠀⠀⠙⢿⠁⠀⠀⠀⠀⠀⠀⡠⠊⠀⠀⠀⢰⡀⠀⣠⠤⠤⢀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀
⠀⣼⣿⠋⠀⠀⠀⠀⢀⣴⣿⠃⠀⠀⠀⢀⣴⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣾⣿⣶⣶⣶⣿⡍⠉⠹⣷⣾⣿⣷⣢⣄⡀⠀⠀⠀⠀⠀⣀⠄⠊⠀⠀⠀⢀⠠⣐⣩⠖⠁⠀⠀⠀⠉⠢⣄⠀⠀⠀⠀⠀⠀⠀⠀⠀
⢰⣿⠃⠀⠀⠀⢀⣴⣿⣿⠃⠀⠀⠀⣠⣾⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣽⣿⣿⣿⣿⣿⣇⠀⠀⠈⠻⡹⡟⢿⠛⠈⠁⠒⢆⣴⣾⡕⠀⡀⠠⣤⣶⣷⣿⣿⡏⠀⠀⠀⠀⠀⠀⠀⠈⠳⡄⠀⠀⠀⠀⠀⠀⠀
⡿⠁⠀⠀⠀⣠⣾⣿⣿⠃⠀⠀⢀⣴⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⠀⠀⠀⠀⢱⡇⠘⠀⠀⢀⣴⣿⣿⠛⠁⠀⠀⠈⠁⠀⠉⠙⡿⠀⠀⠀⠀⠀⠀⠀⠀⢀⡄⠈⢢⠀⠀⠀⠀⠀⠀
⠁⠀⠀⣠⣾⣿⣿⣿⠏⠀⢀⣴⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⠀⠀⠀⠀⣸⣁⣀⣤⣶⣿⣿⣿⣿⣿⡿⠋⠀⠀⠀⠀⢀⡼⠁⠀⠀⠀⠀⠀⠀⢰⡘⢮⠙⣆⠀⠱⡄⠀⠀⠀⠀
⠀⢀⣼⣿⣿⣿⣿⡏⠀⣰⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⡿⠿⠤⠲⠒⠊⣩⣴⣿⣿⣿⣿⣿⣿⣿⣯⠀⠀⠀⠀⠀⣠⡟⠁⠀⠀⠀⠀⠀⠀⠀⠀⠙⠴⠃⠀⠀⠀⠘⢄⠀⠀⠀
⣴⣿⣿⣿⣿⣿⣿⠀⣼⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⠟⠁⠀⠀⣰⣾⣿⣿⣿⣷⣿⣿⣿⣿⠿⢿⡇⠀⠀⢀⣴⠋⡇⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣀⣤⣶⣾⣢⠀⠀
⣿⣿⣿⣿⣿⣿⣿⢰⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⠉⠀⠀⠀⡠⠾⠟⢛⠟⠙⢋⣽⠿⠋⠉⠀⠀⠀⣱⡀⢀⡾⠁⠀⡇⠀⠀⠀⠖⠒⢄⠀⠀⠀⠀⠀⡀⢀⣽⣿⣿⣿⣿⣷⣄
⣿⣿⣿⣿⣿⣿⡇⣾⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⡄⠄⠀⠀⠀⠀⠔⠁⠀⠀⠊⣡⠒⠀⠠⠀⠀⠀⠌⣱⡗⠁⣶⡇⢳⠀⠀⠀⠀⠀⠀⠳⡀⠀⢠⣮⣶⣿⣿⣿⣿⣿⣿⣿⣿
⣿⣿⣿⣿⣿⣿⣧⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣟⡁⠀⠀⠀⠀⣠⠊⠀⢠⠊⠀⢐⣥⡀⢀⠀⡀⢠⣌⡄⠉⠀⠀⢻⡃⢘⡆⠀⠀⠀⠀⠀⠐⢨⢶⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿
⣿⣿⣿⣿⣿⣿⣿⢿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿⣆⠀⡠⠊⠁⠀⣰⣥⣶⣤⡀⠈⠉⠋⠚⠾⢽⡿⠍⠀⠀⠀⠀⢳⡐⣈⠂⠀⠀⠀⠀⠀⡇⠀⠻⣿⣿⣿⣿⣿⣿⣿⣿⣿⣿
⣿⣿⣿⣿⣿⣿⣿⡆⠉⠙⠛⠻⠿⠿⠿⣿⣿⡿⠿⠿⠿⠿⠿⠿⠿⠿⠟⠛⠛⠛⠛⠉⣻⣿⣿⡿⠎⠀⠀⠀⣰⣿⣿⣿⣿⣷⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠹⣷⢢⠀⠀⠀⢀⡜⠀⠀⠀⠈⢿⣿⣿⣿⣿⣿⣿⣿⣿
⣿⣿⣿⣿⣿⣿⣿⣿⣆⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⣿⣿⠟⠁⠀⠀⠀⣼⣿⣿⣿⣿⣿⣿⡇⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠀⠈⣦⣤⣀⣠⡾⠀⠀⠀⠀⠀⠀⠙⢿⣿⣿⣿⣿⣿⣿'
figlet "Konnichiwa $USER!"

# End of file
